# supervisord.conf
[supervisord]
nodaemon=true

[program:ingress]
command=env/bin/python main.py
stdout_logfile=ingress.out
stderr_logfile=ingress.err
