# supervisord.conf
[supervisord]
nodaemon=true
logfile=log/supervisord.log

[program:ingress]
command=env/bin/python -u main.py
stdout_logfile=log/ingress.out
stderr_logfile=log/ingress.err
